# scripts/com.jobpulse.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jobpulse.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/anaconda3/bin/python</string>
        <string>-m</string>
        <string>jobpulse.runner</string>
        <string>daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yashbishnoi/Downloads/multi_agent_patterns</string>

    <!-- Start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if it crashes (after 5 seconds) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <!-- Log stdout/stderr -->
    <key>StandardOutPath</key>
    <string>/Users/yashbishnoi/Downloads/multi_agent_patterns/logs/daemon-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/yashbishnoi/Downloads/multi_agent_patterns/logs/daemon-stderr.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
